# Library/LaunchAgents/io.github.chezmoi.auto-readd.plist.tmpl
{{- /* LaunchAgent to auto re-add selected dotfiles into chezmoi source */ -}}
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>io.github.chezmoi.auto-readd</string>
    <key>ProgramArguments</key>
    <array>
      <string>/bin/zsh</string>
      <string>-lc</string>
      <string>"$HOME/.local/bin/chezmoi-readd.sh"</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>5</integer>
    <key>WatchPaths</key>
    <array>
      <string>{{ .chezmoi.homeDir }}/.zprofile</string>
      <string>{{ .chezmoi.homeDir }}/.zshrc</string>
    </array>
    <key>StandardOutPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/chezmoi-auto-readd.log</string>
    <key>StandardErrorPath</key>
    <string>{{ .chezmoi.homeDir }}/Library/Logs/chezmoi-auto-readd.err.log</string>
  </dict>
  </plist>

